# rpm/dnsproxy.service
[Unit]
Description=DNS Proxy
After=syslog.target
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/sbin/dnsproxy server
TimeoutSec=5000
Restart=on-failure

[Install]
WantedBy=multi-user.target
